# client/config/rmm-agent.service
[Unit]
Description=RMM Agent - Remote Monitoring and Management Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rmm-agent
ExecStart=/usr/bin/node src/agent.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
